# debian/ganeti.postinst.in
#!/bin/sh
# postinst script for ganeti

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
#GANETI_USERS#

	update_symlinks=0
	# Handle symlinks...
	# ... if we are upgrading from pre-2.10 versions
	if [ -n "$2" ] && dpkg --compare-versions "$2" lt 2.10.0; then
		update_symlinks=1

	# ... or if we are not part of a cluster
	elif ! /usr/lib/ganeti/@VER@/usr/lib/ganeti/daemon-util check-config 2>/dev/null; then
		update_symlinks=1

	# ... or if something is broken
	elif [ ! -e /etc/ganeti/lib -o ! -e /etc/ganeti/share ]; then
		update_symlinks=1
	fi

	if [ "$update_symlinks" = 1 ]; then
		for target in lib share; do
			if [ "$(readlink /etc/ganeti/${target})" != "/usr/${target}/ganeti/@VER@" ]; then
				echo "Linking /etc/ganeti/${target} to /usr/${target}/ganeti/@VER@"
				ln -snf /usr/${target}/ganeti/@VER@ /etc/ganeti/${target}
			fi
		done
	fi

	# Restart the service if we have either updated the symlinks, or
	# shipped the same minor version. We only stop the service here and
	# let the debhelper snippet below start it again.
	if [ "$update_symlinks" = 1 -o "${2%.*}" = "@VER@" ]; then
		if [ -x "/etc/init.d/ganeti" ]; then
			invoke-rc.d ganeti stop || true
		fi
	fi

	# Work around bugs of older (pre 2.9) versions. In this case all
	# daemons have already been stopped by the old package's prerm.
	if [ -n "$2" ] && dpkg --compare-versions "$2" lt 2.9.0; then
		# Since 2.9.0, mond runs as root, so remove the old gnt-mond user (that
		# was also created during 2.9.0~rcX.
		if getent passwd gnt-mond >/dev/null; then
			find /var/log/ganeti -user gnt-mond -exec chown root:root {} \;
			deluser --system gnt-mond || true
		fi

		# Remove stale PID files
		# Ganeti versions prior to 2.8.0~rc4 leave Haskell daemon pidfiles
		# behind. If these are root-owned, starting the daemons as unprivileged
		# users will fail.

		# All daemons have been stopped, so it's safe to remove
		# their PID files.
		for daemon in confd mond luxid; do
			rm -f /var/run/ganeti/ganeti-${daemon}.pid
		done

		# Also remove luxid's socket
		rm -f /var/run/ganeti/socket/ganeti-query
	fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# Fix jobqueue archive permissions
# This may take a while, so we only do this if necessary
find /var/lib/ganeti/queue/archive -mindepth 1 -maxdepth 1 \
		-type d -not -user "@GNTMASTERUSER@" 2>/dev/null | while read dirname; do
	/usr/lib/ganeti/ensure-dirs --full-run
	break
done

exit 0
